# main/postmarketos-ui-retroarch/postmarketos-ui-retroarch-openrc.post-install
#!/bin/sh

rc-update add elogind default
rc-update add retroarch default
