# extra-repos/systemd/systemd-services/user-mmsd-tng.service
[Unit]
Description=Multimedia Messaging Service Daemon

[Service]
Type=simple
ExecStart=/usr/bin/mmsdtng
Restart=on-failure

[Install]
WantedBy=default.target
