# extra-repos/systemd/systemd-services/user-plasma-kded6.service
[Unit]
Description=KDE Daemon 6
PartOf=graphical-session.target
After=plasma-kcminit.service

[Service]
ExecStart=/usr/bin/kded6
BusName=org.kde.kded6
Slice=session.slice
TimeoutStopSec=5sec
