# device/community/device-nokia-n900/udev/15-bq27k.rules
KERNEL=="bq27200-0", SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/usr/bin/bq27k-monitor"
