# modem/unisoc-nvitem-daemon/unisoc-nvitem-daemon.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="unisoc-nvitem-daemon"
description="Unisoc NV item daemon"

command="/usr/sbin/unisoc-nvitem-daemon"
command_args="/dev/wwan0sprdnv0"
