# device/community/wl1251-cal/wl1251-cal.initd
#!/sbin/openrc-run
description="Update wl1251 calibration data"

depend()
{
	before modules
	after localmount
}

start()
{
	calibrate-wl1251
	eend $?
}
