# device/community/wl1251-cal/wl1251-cal.initd
#!/sbin/openrc-run
description="Update wl1251 calibration data"

depend()
{
	need localmount
}

start()
{
	calibrate-wl1251
	# Probe wl1251 once the nvs is ready
	modprobe wl1251_spi
	eend $?
}
