# temp/81voltd/81voltd.service
[Unit]
Description=81voltd
After=ModemManager.target

[Service]
ExecStart=/usr/bin/81voltd

[Install]
WantedBy=multi-user.target

